# debian/postinst
#!/bin/bash
# postinst script for bigedit

set -e

ENABLE_SYSTEMD="${BIGEDIT_ENABLE_SYSTEMD:-yes}"

case "$1" in
    configure)
        # Notify user about systemd service
        echo ""
        echo "╔══════════════════════════════════════════════════════════════════╗"
        echo "║                    BigEdit Installation Complete                  ║"
        echo "╠══════════════════════════════════════════════════════════════════╣"
        echo "║                                                                   ║"
        echo "║  BigEdit includes a systemd user service that automatically      ║"
        echo "║  mounts FUSE views for files being edited, so other programs     ║"
        echo "║  can see your changes even before you write them to disk.        ║"
        echo "║                                                                   ║"
        echo "╚══════════════════════════════════════════════════════════════════╝"
        echo ""
        
        if [ "$ENABLE_SYSTEMD" = "yes" ]; then
            echo "The bigedit-watcher service is enabled by default for new user sessions."
            echo ""
            echo "To start it now without logging out, run:"
            echo ""
            echo "    systemctl --user daemon-reload && systemctl --user enable --now bigedit-watcher.service"
            echo ""
            echo "To disable it:"
            echo ""
            echo "    systemctl --user disable --now bigedit-watcher.service"
            echo ""
        else
            echo "Systemd service installation was skipped."
            echo "To enable it manually, run:"
            echo ""
            echo "    systemctl --user enable --now bigedit-watcher.service"
            echo ""
        fi
        
        echo "Usage:"
        echo "    bigedit <filename>     Edit a file"
        echo ""
        echo "Keybindings:"
        echo "    Ctrl+O     Save (journal mode)"
        echo "    Ctrl+J     Compact (full file rewrite)"
        echo "    Ctrl+T     Toggle FUSE mode"
        echo "    Ctrl+X     Exit"
        echo "    Ctrl+G     Help"
        echo ""
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
